# systemd/purlovia-runner.timer
[Unit]
Description=Run Purlovia every day

[Timer]
Unit=purlovia-runner.service
OnCalendar=*-*-* 06,18:02:00
Persistent=true

[Install]
WantedBy=timers.target
